# packages/red-fab-deploy2/red-fab-deploy2-0.2.4.tar.gz/red-fab-deploy2-0.2.4/fab_deploy2/default-configs/templates/base/gunicorn/supervisor_gunicorn.conf
[program:{{ gunicorn.gunicorn_name }}]
{% block config %}
command=sh {{ gunicorn.start_script }}
directory= {{ code_path }}{{ gunicorn.project_path }}
user= {{ gunicorn.user }}
autostart=true
autorestart=true
redirect_stderr=True
{% endblock %}
