# packages/red-fab-deploy2/red-fab-deploy2-0.2.4.tar.gz/red-fab-deploy2-0.2.4/fab_deploy2/default-configs/templates/base/gunicorn/upstart_gunicorn.conf
description "gunicorn server"
author "RED"

start on runlevel [2345]
stop on runlevel [06]

respawn
respawn limit 10 5

script
{{ python.location }}bin/gunicorn_django -c {{ gunicorn.config_location }}
end script
