# packages/Flask-Project/Flask-Project-0.1.0.tar.gz/Flask-Project-0.1.0/flask_project/templates/default/deploy/supervisor.conf
[program:#{project}]
command=/var/www/%(PROJECT_NAME)s/run.sh %(HOST)s %(PORT)d
directory=/var/www/%(project)s
user=root
autostart=true
autorestart=true
environment = SYSTEM_CONFIG="production"