# packages/red-fab-deploy2/red-fab-deploy2-0.2.3.tar.gz/red-fab-deploy2-0.2.3/fab_deploy2/default-configs/templates/base/celery/supervisor_celeryd.conf
[program:{{ celeryd.name }}]
{% block config %}
command=sh {{ celeryd.start_script }}
directory= {{ code_path }}/project
user= {{ celeryd.user }}
autostart=true
autorestart=true
redirect_stderr=True
{% endblock %}

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=998
